# launchd/com.agentflow.computer-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.agentflow.computer-mcp</string>
  <key>ProgramArguments</key>
  <array>
    <string>__ENTRYPOINT__</string>
    <string>--mode</string>
    <string>ws</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>__HOME__/Library/Logs/agentflow-computer-mcp.log</string>
  <key>StandardErrorPath</key><string>__HOME__/Library/Logs/agentflow-computer-mcp.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>__USER_BIN__:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
